# docker/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:teleirc]
command=teleirc

[program:nginx]
command=nginx -g "daemon off;"